# configs/cliproxyapi/Library/LaunchAgents/com.cliproxyapi.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cliproxyapi</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/opt/cliproxyapi/bin/cliproxyapi</string>
        <string>-config</string>
        <string>/Users/ice/.config/cliproxyapi/config.yaml</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/cliproxyapi.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/cliproxyapi.err</string>
</dict>
</plist>
